# Library/LaunchAgents/me.diti.scripts.brewupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>me.diti.scripts.brewupdate</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/brew</string><!-- FIXME: Won’t work at 42, because ~/.brew/bin/brew instead! -->
    <string>update</string>
  </array>

  <key>LowPriorityIO</key>
  <true/>

  <key>StartInterval</key>
  <integer>21600</integer><!-- 6 hours -->

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/brewupdate-out.txt</string>

  <key>StandardErrorPath</key>
  <string>/tmp/brewupdate-err.txt</string>
</dict>
</plist>
